# spotify/rootfs/etc/services.d/spotifyd/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Spotify Connect
# Runs librespot
# ==============================================================================
declare -a options
declare bitrate
declare name
declare password
declare username
declare autoplay

bashio::log.info 'Starting the Spotify daemon...'

# Bitrate
options+=(--bitrate $(bashio::config 'bitrate'))

# Device name
name=$(bashio::config 'name')
options+=(--name "${name}")

# Username / password
if bashio::config.has_value 'username'; then
    bashio::config.require.password

    username=$(bashio::config 'username')
    password=$(bashio::config 'password')

    options+=(--username "${username}")
    options+=(--password "${password}")
fi

# Autoplay
if bashio::config.true 'autoplay'; then
  options+=(--autoplay)
fi

# Save writes
options+=(--disable-audio-cache)

# Are we running in debug mode?
if bashio::debug; then
  options+=(--verbose)
fi

# Run librespot
exec librespot "${options[@]}"
