# spotify/rootfs/etc/services.d/spotifyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Spotify Connect
# Runs spotifyd
# ==============================================================================
declare -a options

bashio::log.info 'Starting the Spotify daemon...'

# Run spotifyd in the foreground
options+=(--no-daemon)

# Are we running in debug mode?
if bashio::debug; then
  options+=(--verbose)
fi

# Run spotifyd
exec spotifyd "${options[@]}"
